# templates/mattermost.service.tmpl
[Unit]
Description=Mattermost

[Service]
Type=simple
User=mattermost
Group=mattermost
WorkingDirectory=/srv/mattermost
ExecStart=/srv/mattermost/bin/platform
Restart=always
RestartSec=1
